# services/lektrix.update.service
# This service is for updating the repository clone

[Unit]
Description=lektrix update service
After=multi-user.target
Wants=multi-user.target
Wants=lektrix.update.timer

[Service]
Type=oneshot
User=pi
ExecStart=/home/pi/lektrix/lektrix --systemd
TimeoutStopSec=900
