# services/lektrix.prices.timer
# This timer is for triggering the daily collection of hourly electricity prices

[Unit]
Description=collecting electricity prices (timer)

[Timer]
OnCalendar=*-*-* 12,13,14,15,16,17,18,23:03
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
